# shadow.timer
[Unit]
Description=Daily verification of password and group files

[Timer]
OnCalendar=daily
AccuracySec=12h
RandomizedDelaySec=5m
Persistent=true
